# launchd/com.assistant.vault-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.assistant.vault-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/joshlevine/.bun/bin/bun</string>
        <string>run</string>
        <string>src/scheduler/runner.ts</string>
        <string>vault-sync</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joshlevine/src/tries/2025-12-25-agent-orchestrator/personal-assistant</string>

    <!-- Run every 4 hours -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>/Users/joshlevine/.assistant/logs/vault-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/joshlevine/.assistant/logs/vault-sync.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/joshlevine/.bun/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
